# templates/opendkim.service.j2
# {{ ansible_managed }}

[Service]
PIDFile=
User={{ opendkim_user }}
Group={{ opendkim_group }}
